# build-images/rl-image/default.ini
[program:cron]
command=/supervisor/cron.sh
user=root

[program:rsyslog]
command=/supervisor/rsyslog.sh
user=root

[program:sshd]
command=/supervisor/sshd.sh
user=root

[program:sysstat]
startsecs=0
autorestart=false
command=/usr/lib64/sa/sa1 --boot
user=root
