# deploy/zimbra/extras/zmstat-cleanup.cron
#
# zmstat_cleanup
#
15 0 * * 7 /opt/zimbra/libexec/zmstat-cleanup -k 30
